# devel/etc/catkin/profile.d/python_orocos_kdl_site_packages.zsh
export PYTHONPATH=/home/seungbeom/franka_ws/devel/lib/python2.7/dist-packages:$PYTHONPATH
